# src/etc/init.d/sdc-init
#!/sbin/openrc-run
#
# Copyright 2020 Joyent, Inc.
#

description="Initializes the zone"

start()
{
	ebegin "Initializing Zone"
	/lib/smartdc/joyent_rc.local
	eend $?
}
