# src/etc/profile.d/native_manpath.sh
# Include man pages for mdata-* binaries in /native
export MANPATH=$MANPATH:/native/usr/share/man
