# src/etc/profile.d/native_manpath.sh
# Include man pages for wrapper scripts and symlinks
# that reference binaries in /native
export MANPATH=$MANPATH:/native/usr/share/man
